# manager/com.openclaw.trading.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.trading</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/roberto/Documents/OpenClawAgents/manager/process_manager_v2.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/roberto/Documents/OpenClawAgents/logs/manager.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/roberto/Documents/OpenClawAgents/logs/manager.error</string>
    <key>WorkingDirectory</key>
    <string>/Users/roberto/Documents/OpenClawAgents</string>
</dict>
</plist>
